# scripts/launchd/useorgx.agents-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>useorgx.agents-refresh</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>node</string>
      <string>__REPO_ROOT__/scripts/daily-agents-refresh.mjs</string>
      <string>--apply=true</string>
      <string>--scope=code</string>
      <string>--engine=codex</string>
    </array>

    <!-- Daily at 02:17 local time -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>2</integer>
      <key>Minute</key>
      <integer>17</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/.codex/reports/agents-refresh/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.codex/reports/agents-refresh/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
